# ircddbgatewayd.timer
[Timer]
OnStartupSec=31

[Install]
WantedBy=multi-user.target
