# service/ysfgateway.timer
[Timer]
OnStartupSec=35

[Install]
WantedBy=multi-user.target
